# scripts/gpu-powercap.service
[Unit]
Description=Cap NVIDIA GPU power + lock clocks (PSU stability)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/nvidia-smi -pm 1
ExecStart=/usr/bin/nvidia-smi -pl 200
ExecStart=/usr/bin/nvidia-smi -lgc 300,1290
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
